# config/dansguardian_autostart.conf
description "dansguardian"
author "Dry1d"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
	/etc/init.d/dansguardian start
end script
