# config/squid3_autostart.conf
description "squid3"
author "Dry1d"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
	/etc/init.d/squid3 start
end script
